# lab4/ansible/roles/webapp/templates/mywebapp.service.j2
[Unit]
Description=Simple Inventory Web Application
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory=/opt/{{ app_name }}
ExecStart=/opt/{{ app_name }}/venv/bin/python /opt/{{ app_name }}/app.py \
    --host 127.0.0.1 \
    --port {{ app_port }} \
    --db-host {{ db_host }} \
    --db-port {{ db_port }} \
    --db-user {{ db_user }} \
    --db-password {{ db_password }} \
    --db-name {{ db_name }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
